# daemon/stalking.timer
[Unit]
Description=stalking

[Timer]
OnUnitActiveSec=1000s
OnBootSec=1000s

[Install]
WantedBy=default.target