# daemon/stalking.timer
[Unit]
Description=stalking

[Timer]
OnUnitActiveSec=300s
OnBootSec=300s

[Install]
WantedBy=default.target